# nagios/sshd.ini
;sshd 
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/sshd/%(program_name)s.error.log
stderr_logfile=/var/log/sshd/%(program_name)s.log
autorestart=true
